# crontab/getRFIstats.sh
#!/bin/sh
cd /var/www/html/prod.epicore.org/scripts/
php downloadEventStats.php
